# roles/spark_observability/templates/ic-probe-textfile.service.j2
[Unit]
Description=Write Spark interconnect probe metrics to node_exporter textfile dir

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ic-probe-textfile.sh
